# roles/mdns-socat-tunnel/templates/socat-mdns-receiver.service.j2
[Unit]
Description=socat mDNS receiver service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
RuntimeMaxSec=10800
User=root
ExecStart=/usr/local/bin/socat_mdns_receiver.sh

[Install]
WantedBy=multi-user.target
